# conf/supervisor.d/sssd.ini
[program:sssd]
autostart=true
command=/usr/sbin/sssd -i -d 4
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
priority=50
